# deploy/systemd/gback.service
[Unit]
Description=gback Node.js API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/gback
EnvironmentFile=/etc/gback.env
Environment=NODE_ENV=production
LimitNOFILE=100000
ExecStart=/usr/bin/node /opt/gback/index.js
Restart=always
RestartSec=3
User=www-data
Group=www-data

# Hardening (safe defaults)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
